# board/myna-player-odyssey/rootfs_overlay/etc/systemd/system/rauc-failure.service
[Unit]
Description=Unit triggered on rauc-marking-good failing
SuccessAction=reboot

[Service]
Type=oneshot
ExecStart=/bin/true
