# board/Phoenix/rootfs_overlay/etc/systemd/system/rngd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/sbin/rngd -f -x jitter
